# scripts/systemd/dp-tasks.timer
[Unit]
Description=Periodically launch tasks

[Timer]
OnBootSec=30s
OnUnitActiveSec=15min
